# files/docker-registry.conf
description "Docker Registry"
version "0.6.3"
author "Abhinav Pobbati"

start on filesystem and started docker
stop on runlevel [2345]
respawn
script
        /usr/bin/docker start registry
end script

